# benchmarks-website/ops/systemd/vortex-bench-backup.timer
# SPDX-License-Identifier: Apache-2.0
# SPDX-FileCopyrightText: Copyright the Vortex contributors
#
# Hourly snapshot timer. Persistent=true means a missed hour catches up
# after a reboot.

[Unit]
Description=Vortex bench v3 hourly DB snapshot timer
Requires=vortex-bench-backup.service

[Timer]
# Top of every hour, UTC.
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=2min
Unit=vortex-bench-backup.service

[Install]
WantedBy=timers.target
